# rgb/.config/systemd/user/lianli-daemon.service
[Unit]
Description=Lian Li Device Daemon

[Service]
Type=simple
ExecStart=/usr/bin/lianli-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
