# rgb/.config/systemd/user/rgb-off.service
[Unit]
Description=Turn off RAM and Lian Li fan RGB
After=lianli-daemon.service
Wants=lianli-daemon.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rgb-off
RemainAfterExit=no

[Install]
WantedBy=default.target
